# templates/node_exporter.service.erb
[Unit]
Description=Node Exporter <%= node['node_exporter']['version'] %>

[Service]
Type=simple
User=<%= node['node_exporter']['user'] %>
Group=<%= node['node_exporter']['group'] %>
PIDFile=<%= node['node_exporter']['path']['pid'] %>
ExecStart=<%= node['node_exporter']['path']['bin'] %> <%= generate_node_exporter_config %>
Restart=on-failure

[Install]
WantedBy=multi-user.target